# rescue/tree/etc/inittab
#Inittab for Mandrake Rescue Disk
#Author: Daouda Lo <daouda@mandrakesoft.com>
id:3:initdefault:

si::sysinit:/etc/rc.sysinit
re:6:wait:/etc/rc.reboot
ca:12345:ctrlaltdel:/sbin/fakeshutdown

1:12345:respawn:/sbin/mingetty --autologin root tty1
2:23:respawn:/sbin/mingetty tty2
3:23:respawn:/sbin/mingetty tty3
4:23:respawn:/sbin/mingetty tty4
5:23:respawn:/sbin/mingetty tty5
6:23:respawn:/sbin/mingetty tty6
